# roles/rancher_k8s/templates/rancher.service.j2
[Unit]
Description=Rancher Kubernetes Management Platform
Documentation=https://rancher.com/docs/
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/rancher
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down
TimeoutStartSec=180
TimeoutStopSec=120
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
