# roles/redis_cache/templates/redis-sentinel.service.j2
[Unit]
Description=Redis Sentinel
After=network.target redis.service
Documentation=https://redis.io/documentation

[Service]
Type=notify
ExecStart=/usr/bin/redis-sentinel {{ redis_config_dir }}/sentinel.conf --supervised systemd
ExecStop=/usr/bin/redis-cli -p {{ redis_sentinel_port }} shutdown
User={{ redis_user }}
Group={{ redis_user }}
RuntimeDirectory=redis-sentinel
RuntimeDirectoryMode=0755
LimitNOFILE=65536
TimeoutStopSec=60
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
